# configs/local/Library/LaunchAgents/local.cliproxyapi.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>local.cliproxyapi</string>
    
    <key>ProgramArguments</key>
    <array>
      <string>/bin/sh</string>
      <string>-c</string>
      <string>PATH=/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:$HOME/.nix-profile/bin exec /opt/homebrew/bin/cliproxyapi --config $HOME/.config/cliproxyapi/config.yaml >> /tmp/cliproxyapi.log 2>&1</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
      <key>SuccessfulExit</key>
      <false/>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/thanhhaikhong/.nix-profile/bin</string>
      <key>HOME</key>
      <string>/Users/thanhhaikhong</string>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/tmp/cliproxyapi.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/cliproxyapi.log</string>
  </dict>
</plist>